# launchd/com.galew.speak-to-user.stable.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.galew.speak-to-user.stable</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec "$HOME/Workspace/speak-to-gale/scripts/run_service.sh" &gt;&gt;"$HOME/Library/Logs/speak-to-user.stable.stdout.log" 2&gt;&gt;"$HOME/Library/Logs/speak-to-user.stable.stderr.log"</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SPEAK_TO_USER_HOST</key>
        <string>127.0.0.1</string>
        <key>SPEAK_TO_USER_PORT</key>
        <string>8765</string>
        <key>SPEAK_TO_USER_PLAYBACK_BACKEND</key>
        <string>wavbuffer</string>
        <key>SPEAK_TO_USER_WAVBUFFER_QUEUE_DEPTH</key>
        <string>8</string>
        <key>SPEAK_TO_USER_WAVBUFFER_PREROLL_MODE</key>
        <string>buffers</string>
        <key>SPEAK_TO_USER_PLAYBACK_PREROLL_CHUNKS</key>
        <string>3</string>
        <key>SPEAK_TO_USER_PLAYBACK_PREROLL_SECONDS</key>
        <string>5.0</string>
        <key>SPEAK_TO_USER_TTS_CHUNK_MAX_CHARS</key>
        <string>240</string>
        <key>SPEAK_TO_USER_TTS_MAX_NEW_TOKENS</key>
        <string>384</string>
        <key>SPEAK_TO_USER_TTS_MAX_CHUNK_SYNTH_SECONDS</key>
        <string>30.0</string>
        <key>SPEAK_TO_USER_TTS_MAX_CHUNK_AUDIO_SECONDS</key>
        <string>20.0</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
