# launchd/com.galew.speak-to-user.stable.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.galew.speak-to-user.stable</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/galew/Workspace/speak-to-gale/scripts/run_service.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SPEAK_TO_USER_HOST</key>
        <string>127.0.0.1</string>
        <key>SPEAK_TO_USER_PORT</key>
        <string>8765</string>
        <key>SPEAK_TO_USER_PLAYBACK_BACKEND</key>
        <string>wavbuffer</string>
        <key>SPEAK_TO_USER_PLAYBACK_PREROLL_SECONDS</key>
        <string>5.0</string>
        <key>SPEAK_TO_USER_WAVBUFFER_BINARY_PATH</key>
        <string>/Users/galew/Workspace/swiftly-play/.build/release/wavbuffer</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/galew/Library/Logs/speak-to-user.stable.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/galew/Library/Logs/speak-to-user.stable.stderr.log</string>
</dict>
</plist>
