# launchd/com.galew.speak-to-user.dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.galew.speak-to-user.dev</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec "$HOME/Workspace/speak-to-user/scripts/run_service.sh" &gt;&gt;"$HOME/Library/Logs/speak-to-user.dev.stdout.log" 2&gt;&gt;"$HOME/Library/Logs/speak-to-user.dev.stderr.log"</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SPEAK_TO_USER_HOST</key>
        <string>127.0.0.1</string>
        <key>SPEAK_TO_USER_PORT</key>
        <string>8766</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
